# crontab
3 2 1 * * python C:\Users\icarl\OneDrive - Oliver Maura Consultores, S.L\Desktop\Formaciones\software_engineering\software_engineering_for_data_scientists\ch11\pipeline_for_scheduling.py
